# docker/supervisord-render.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/tmp/supervisord.pid
environment=HOME="/tmp"

[program:api]
command=python -m uvicorn src.aicbc.main:app --host 0.0.0.0 --port 8000
directory=/app
user=aicbc
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=celery -A src.aicbc.analysis.tasks worker --loglevel=info --concurrency=1 -Q analysis,celery
directory=/app
user=aicbc
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:beat]
command=celery -A src.aicbc.analysis.tasks beat --loglevel=info --scheduler celery.beat.PersistentScheduler -s /tmp/celerybeat-schedule
directory=/app
user=aicbc
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
